# check-engine/check-engine.initd
#!/sbin/openrc-run
name="check-engine"
description="EasyGuard System Watcher"
command="check-engine"
command_args="-l"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
