# failoverd/failoverd.initd
#!/sbin/openrc-run
name="failoverd"
description="Failover Daemon"
command="failoverd"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	provide net
}
